# scraper.cron
# Run the lunch menu scraper every night at 10 PM
0 22 * * * root cd /work && /usr/local/bin/python /work/scripts/scrape_weekly_menu.py >> /var/log/cron.log 2>&1

